# deploy/systemd/admind.service
[Unit]
Description=OpenTrusty Admin Plane (admind)
After=network.target postgresql.service

[Service]
Type=simple
User=opentrusty
Group=opentrusty
EnvironmentFile=/etc/opentrusty/admind.env
ExecStart=/usr/local/bin/admind
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
